# washdb-bot/systemd/yp-state-workers.service
[Unit]
Description=YP State Worker Pool (Optimized with Browser Pooling)
After=network.target postgresql.service

[Service]
Type=simple
User=rivercityscrape
Group=rivercityscrape

# Load environment variables from config
EnvironmentFile=/etc/washdb-bot/washdb-bot.env
EnvironmentFile=${WASHDB_HOME}/.env

WorkingDirectory=${WASHDB_HOME}
Environment="PATH=${WASHDB_VENV}/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Run the optimized state worker pool
ExecStart=${WASHDB_VENV}/bin/python -m scrape_yp.state_worker_pool

# Restart policy
Restart=always
RestartSec=10

# Logging (both journal and file for GUI live output)
StandardOutput=append:${WASHDB_LOGS}/yp_workers.log
StandardError=append:${WASHDB_LOGS}/yp_workers.log
SyslogIdentifier=yp-state-workers

[Install]
WantedBy=multi-user.target
