# washdb-bot/services/llm-scheduler.service
[Unit]
Description=LLM Service Scheduler - Manages verification and standardization
After=network.target ollama.service postgresql.service
Wants=ollama.service

[Service]
Type=simple
User=rivercityscrape
Group=rivercityscrape
WorkingDirectory=/home/rivercityscrape/URL-Scrape-Bot/washdb-bot
Environment="PATH=/home/rivercityscrape/URL-Scrape-Bot/washdb-bot/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DISPLAY=:99"
ExecStart=/home/rivercityscrape/URL-Scrape-Bot/washdb-bot/venv/bin/python scripts/llm_scheduler.py --verification-start 6 --verification-end 18
Restart=always
RestartSec=60
StandardOutput=append:/home/rivercityscrape/URL-Scrape-Bot/washdb-bot/logs/llm_scheduler.log
StandardError=append:/home/rivercityscrape/URL-Scrape-Bot/washdb-bot/logs/llm_scheduler.log

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
