# washdb-bot/washbot-yp-scraper.service
[Unit]
Description=Washbot YP Scraper Service (24/7)
Documentation=https://github.com/rivercityscrape/washdb-bot
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=rivercityscrape
Group=rivercityscrape

# Load environment variables from config
EnvironmentFile=/etc/washdb-bot/washdb-bot.env
EnvironmentFile=${WASHDB_HOME}/.env

WorkingDirectory=${WASHDB_HOME}
Environment=PYTHONUNBUFFERED=1
ExecStart=${WASHDB_VENV}/bin/python scripts/yp_continuous.py

# Restart configuration
Restart=always
RestartSec=30

# Logging
StandardOutput=append:${WASHDB_LOGS}/yp_service.log
StandardError=append:${WASHDB_LOGS}/yp_service_error.log

# Resource limits
MemoryMax=4G
CPUQuota=50%

# Security
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=${WASHDB_HOME}
PrivateTmp=true

[Install]
WantedBy=multi-user.target
