# washdb-bot/washdb-bot.service
[Unit]
Description=Washdb-Bot NiceGUI Dashboard
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=rivercityscrape
Group=rivercityscrape

# Load environment variables from config
EnvironmentFile=/etc/washdb-bot/washdb-bot.env
EnvironmentFile=${WASHDB_HOME}/.env

WorkingDirectory=${WASHDB_HOME}
Environment="PATH=${WASHDB_VENV}/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=${WASHDB_VENV}/bin/python -m niceui.main

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
